# deploy/donors/qorm-donors.timer
[Unit]
Description=Daily QORM supporters-list refresh

[Timer]
OnCalendar=*-*-* 04:15:00
Persistent=true

[Install]
WantedBy=timers.target
